# systemd/saucery-saucier.timer
[Unit]
Description=Timer for Saucery Saucier grocery shopping

[Timer]
OnUnitInactiveSec=30min
