# systemd/saucery-grocer.timer
[Unit]
Description=Timer for Saucery Grocer shelf stocking

[Timer]
OnUnitInactiveSec=10min
