# systemd/saucery-grocer.timer
[Unit]
Description=Grocery shop for sosreports

[Timer]
OnUnitInactiveSec=300
