# systemd/saucery-saucier.service
[Unit]
Description=Stock Saucery with new sosreports from Grocery

[Service]
ExecStart=saucier -ESs
