# systemd/saucery-grocer.service
[Unit]
Description=Grocery shop for sosreports

[Service]
ExecStart=grocer -s
