# systemd/saucery-grocer-clean.service
[Unit]
Description=Clean empty grocery shelves

[Service]
ExecStart=grocer -C --logname grocer
