# systemd/units/saucery-grocer-stock.timer
[Unit]
Description=Timer for Saucery Grocer shelf stocking

[Timer]
OnUnitInactiveSec=10min
OnStartupSec=10min

[Install]
WantedBy=timers.target
