# systemd/units/saucery-saucier-menu.timer
[Unit]
Description=Timer to rebuild saucery menu

[Timer]
OnUnitInactiveSec=5min
OnStartupSec=5min

[Install]
WantedBy=timers.target
