# systemd/units/saucery-saucier-sear.service
[Unit]
Description=Sear new sosreports

[Service]
ExecStart=saucier -S --logname saucier
User=saucery
