# systemd/units/saucery-grocer-clean.timer
[Unit]
Description=Timer for Saucery Grocer shelf cleaning

[Timer]
OnUnitInactiveSec=1day
OnStartupSec=1day

[Install]
WantedBy=timers.target
