# systemd/units/saucery-saucier-shop.service
[Unit]
Description=Stock Saucery with new sosreports from Grocery

[Service]
ExecStart=saucier -ESMs --logname saucier
User=saucery
