# systemd/units/saucery-saucier-menu.service
[Unit]
Description=Create new saucery menu

[Service]
ExecStart=saucier -M --logname saucier
User=saucery
