# systemd/units/saucery-saucier-extract.timer
[Unit]
Description=Timer to extract new sosreports

[Timer]
OnUnitInactiveSec=15min
OnStartupSec=15min

[Install]
WantedBy=timers.target
