# systemd/units/saucery-saucier-extract.service
[Unit]
Description=Extract new sosreports

[Service]
ExecStart=saucier -E --logname saucier
User=saucery
