# systemd/units/saucery-saucier-sear.timer
[Unit]
Description=Timer to sear new sosreports

[Timer]
OnUnitInactiveSec=15min
OnStartupSec=15min

[Install]
WantedBy=timers.target
