# systemd/units/saucery-grocer-stock.service
[Unit]
Description=Stock Grocery shelves with new deliveries

[Service]
ExecStart=grocer -s --logname grocer
User=saucery
