# systemd/units/saucery-saucier-shop.timer
[Unit]
Description=Timer for Saucery Saucier grocery shopping

[Timer]
OnUnitInactiveSec=15min
OnStartupSec=15min

[Install]
WantedBy=timers.target
